# genomebin/services/rc.d.tmpl
#!/bin/sh
# PROVIDE: {{PRIMAL_NAME}}
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="{{PRIMAL_NAME}}"
rcvar="{{PRIMAL_NAME}}_enable"
desc="{{PRIMAL_DESCRIPTION}}"

load_rc_config $name

: ${{{PRIMAL_NAME}}_enable:="NO"}
: ${{{PRIMAL_NAME}}_user:="{{SERVICE_USER}}"}
: ${{{PRIMAL_NAME}}_group:="{{SERVICE_GROUP}}"}
: ${{{PRIMAL_NAME}}_config:="{{CONFIG_DIR}}/config.toml"}
: ${{{PRIMAL_NAME}}_log:="{{LOG_DIR}}/{{PRIMAL_NAME}}.log"}
: ${{{PRIMAL_NAME}}_pidfile:="/var/run/{{PRIMAL_NAME}}.pid"}

command="{{BINARY_PATH}}"
command_args="{{START_ARGS}}"
pidfile="${{{PRIMAL_NAME}}_pidfile}"

start_precmd="{{PRIMAL_NAME}}_precmd"
stop_postcmd="{{PRIMAL_NAME}}_postcmd"

{{PRIMAL_NAME}}_precmd()
{
    # Ensure log directory exists
    mkdir -p "{{LOG_DIR}}"
    chown "${{{PRIMAL_NAME}}_user}:${{{PRIMAL_NAME}}_group}" "{{LOG_DIR}}"
    
    # Ensure data directory exists
    mkdir -p "{{DATA_DIR}}"
    chown "${{{PRIMAL_NAME}}_user}:${{{PRIMAL_NAME}}_group}" "{{DATA_DIR}}"
    
    # Set environment
    export RUST_LOG="{{RUST_LOG}}"
    export {{PRIMAL_NAME_UPPER}}_CONFIG="${{{PRIMAL_NAME}}_config}"
}

{{PRIMAL_NAME}}_postcmd()
{
    rm -f "$pidfile"
}

run_rc_command "$1"

