# genomebin/services/systemd.service.tmpl
[Unit]
Description={{PRIMAL_NAME}} - {{PRIMAL_DESCRIPTION}}
Documentation=https://docs.ecoprimals.dev/{{PRIMAL_NAME}}/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
ExecStart={{BINARY_PATH}} {{START_ARGS}}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{PRIMAL_NAME}}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{DATA_DIR}} {{CONFIG_DIR}}
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Environment
Environment="RUST_LOG={{RUST_LOG}}"
Environment="{{PRIMAL_NAME_UPPER}}_CONFIG={{CONFIG_DIR}}/config.toml"
EnvironmentFile=-{{CONFIG_DIR}}/environment

[Install]
WantedBy=multi-user.target

